# systemd/click-bridge-heal.service
[Unit]
Description=click-bridge self-heal check (detect + repair + report)

[Service]
Type=oneshot
ExecStart=%h/projects/click-bridge/tools/self-heal.sh
